# vagrant_share/ansible/roles/asl-webserver/templates/update-crl.timer.j2
[Unit]
Description=Update chained CRL

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
