# vagrant_share/ansible/roles/asl-webserver/templates/update-crl.service.j2
[Unit]
Description=Update chained CRL

[Service]
ExecStart=/usr/local/bin/update_crl.sh
