# vagrant_share/ansible/roles/backup-client/templates/borgbackup.service.j2
[Unit]
Description=Borg

[Service]
Type=simple
ExecStart=/etc/borgbackup/run.sh
User=root
Group=root
